# youtube-dl/rootfs/etc/services.d/youtube-dl/finish
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Home Assistant Add-on: youtube-dl
# Take down the S6 supervision tree when youtube-dl fails
# ==============================================================================

if bashio::config.true 'STDIN_mode'; then
    bashio::log.info 'STDIN service was stopped, restarting...'
else
    s6-svscanctl -t /var/run/s6/services
fi
