# youtube-dl/rootfs/etc/services.d/youtube-dl/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Home Assistant Add-on: youtube-dl
# Runs youtube-dl
# ==============================================================================

bashio::log.info 'Starting youtube-dl...'

# Run youtube-dl
if bashio::config.true 'STDIN_mode'; then
    bashio::log.info "Starting in STDIN_mode..."
    while read -r URL_STDIN; do
        bashio::log.info "Received URL: ${URL_STDIN}"
        export ${URL_STDIN}
        exec /usr/bin/youtube-dl.sh
    done
else
    exec /usr/bin/youtube-dl.sh
fi









