# youtube-dl/rootfs/etc/services.d/youtube-dl/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Home Assistant Add-on: youtube-dl
# Runs youtube-dl
# ==============================================================================

bashio::log.info 'Starting youtube-dl...'

# Run youtube-dl
exec /usr/bin/youtube-dl.sh
